# ansible/roles/mailcatcher/files/mailcatcher.conf
[program:mailcatcher]
command=mailcatcher --foreground --http-ip=0.0.0.0
autostart=true
autorestart=true
stderr_logfile=/var/log/mailcatcher.err.log
stdout_logfile=/var/log/mailcatcher.out.log
